# .config/systemd/user/battery-alert.service
[Unit]
Description=Low battery notification (20%% and 5%%)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-alert
