# .config/systemd/user/hypridle.service
[Unit]
Description=Hypridle idle management (AC/battery aware)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hypridle-launch
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
